# bin/updater/systemd/owncloud-news-updater.service
[Unit]
After=default.target

[Service]
Type=simple
User=nobody
ExecStart=/usr/bin/python3 /usr/local/bin/owncloud-news-updater -c /etc/owncloud/news/updater.ini

[Install]
WantedBy=default.target
